# nginxforodoo/crontab
3 16 * * * root certbot renew >> /var/log/nginx/renew.log
